# crontab
18 11,12,14,16,18 * * 1-5 root /app/.venv/bin/python /app/main.py >> /proc/1/fd/1 2>> /proc/1/fd/2
18 14,18 * * 0,6 root /app/.venv/bin/python /app/main.py >> /proc/1/fd/1 2>> /proc/1/fd/2